# dotfiles/windows_only/PowerShell/files/profile.d/00_import.ps1
Import-Module -name Get-ChildItemColor
Import-Module -name PSReadLine
Import-Module -name posh-git
