# dotfiles/windows_only/PowerShell/files/profile.d/02_alias.ps1
Set-Alias ls Get-ChildItemColor -option AllScope -Force
Set-Alias dir Get-ChildItemColor -option AllScope -Force
